# devops/ansible/roles/camwatcher/templates/camwatcher.service.j2
# camwatcher service unit file
# Generated by Ansible template
# Defines camwatcher.py as a systemd service

[Unit]
Description=CamWatcher Service
After=multi-user.target

[Service]
Type=idle
WorkingDirectory={{ camwatcher_install_path }}
User={{ sentinelcam_user }}
Group={{ sentinelcam_group }}
Restart=always
RestartSec=5
ExecStart={{ sentinelcam_services.python_executable }} -u {{ camwatcher_install_path }}/camwatcher.py

[Install]
WantedBy=multi-user.target
