# devops/ansible/roles/datapump/templates/datapump.service.j2
# datapump service unit file
# Generated by Ansible template
# Defines datapump.py as a systemd service

[Unit]
Description=DataPump Service
After=multi-user.target

[Service]
Type=idle
WorkingDirectory={{ datapump_install_path }}
User={{ sentinelcam_user }}
Group={{ sentinelcam_group }}
Restart=always
RestartSec=5
ExecStart={{ sentinelcam_services.python_executable }} -u {{ datapump_install_path }}/datapump.py

[Install]
WantedBy=multi-user.target
