# devops/ansible/roles/imagehub/templates/imagehub.service.j2
# imagehub service unit file
# Generated by Ansible template
# Defines imagehub.py as a systemd service

[Unit]
Description=Imagehub Service
After=multi-user.target

[Service]
Type=idle
WorkingDirectory={{ imagehub_install_path }}/imagehub
User={{ sentinelcam_user }}
Group={{ sentinelcam_group }}
Restart=always
RestartSec=5
ExecStart={{ sentinelcam_services.python_executable }} -u {{ imagehub_install_path }}/imagehub/imagehub.py

[Install]
WantedBy=multi-user.target
