# devops/ansible/roles/watchtower/templates/watchtower.service.j2
# Watchtower service unit file
# Generated by Ansible template
# Defines watchtower.py as a systemd service for GUI display

[Unit]
Description=Sentinelcam Watchtower
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User={{ sentinelcam_user }}
WorkingDirectory={{ watchtower_install_path }}/watchtower
Environment=DISPLAY={{ watchtower_display }}
Environment=XAUTHORITY=/home/{{ sentinelcam_user }}/.Xauthority
Restart=always
RestartSec=5
ExecStart={{ sentinelcam_venv_opencv }}/bin/python -u {{ watchtower_install_path }}/watchtower/watchtower.py

[Install]
WantedBy=graphical.target
