# devops/ansible/roles/imagenode/templates/imagenode.service.j2
# imagenode service unit file
# Generated by Ansible for {{ inventory_hostname }}
# Defines imagenode.py as a systemd service

[Unit]
Description=ImageNode Service for {{ node_name }}
After=multi-user.target network.target

[Service]
Type=idle
WorkingDirectory={{ imagenode_install_path }}/imagenode
User={{ sentinelcam_user }}
Group={{ sentinelcam_group }}
Restart=always
RestartSec=5
ExecStart={{ sentinelcam_services.python_executable }} -u {{ imagenode_install_path }}/imagenode/imagenode.py

# Environment variables
Environment=PYTHONPATH={{ imagenode_install_path }}/imagenode

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=imagenode-{{ node_name }}

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
